# packaging/systemd/shieldbuntu.service
[Unit]
Description=ShieldBuntu hardening service
Documentation=https://github.com/Akshat-Pandey16/ShieldBuntu
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=root
Group=root
ExecStart=/usr/bin/shieldbuntu-server
Restart=on-failure
RestartSec=5s
WorkingDirectory=/var/lib/shieldbuntu

Environment=SHIELDBUNTU_DEV_MODE=false
Environment=SHIELDBUNTU_HOST=127.0.0.1
Environment=SHIELDBUNTU_PORT=8765
Environment=SHIELDBUNTU_LOG_LEVEL=INFO
Environment=SHIELDBUNTU_DATA_DIR=/var/lib/shieldbuntu
Environment=SHIELDBUNTU_ANSIBLE_ROOT=/usr/share/shieldbuntu/ansible

StandardOutput=journal
StandardError=journal
SyslogIdentifier=shieldbuntu

# Note: we deliberately do NOT enable ProtectSystem/ProtectHome/etc.
# This daemon's purpose is to modify /etc and system state — sandboxing
# it the same way as a typical service would defeat the point.

[Install]
WantedBy=multi-user.target
